# rootfs/etc/supervisor.d/carbon-relay-ng.ini
[program:carbon-relay]
command        = /usr/bin/carbon-relay-ng /etc/carbon-relay-ng.ini
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
